# hud35.service
[Unit]
Description=HUD35 Launcher Service
After=network.target
Wants=network.target

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/opt/hud35
ExecStart=/usr/bin/python3 /opt/hud35/launcher.py
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal
[Install]
WantedBy=multi-user.target
